# crontab
 0 22 * * * /home/admbbkrv/www/gr1application/vendor/bin/laminas tokens:expire-soon
